# deploy/systemd/intelligence-backup.timer
[Unit]
Description=Run Intelligence ITA DB backup daily at 03:00 UTC

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
