# deploy/systemd/intelligence-pipeline.service
[Unit]
Description=Intelligence ITA Daily Pipeline
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/intelligence-ita/app
EnvironmentFile=/opt/intelligence-ita/app/.env.production
ExecStart=/usr/bin/docker compose run --rm backend python scripts/daily_pipeline.py
TimeoutStartSec=3600
User=root
