# deploy/systemd/intelligence-backup.service
[Unit]
Description=Intelligence ITA Database Backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/intelligence-ita/app/deploy/backup-db.sh
User=root
